# systemd/btrfs-raid.target
[Unit]
Description=Btrfs raid management
StopWhenUnneeded=true
DefaultDependencies=false
ConditionFileIsExecutable=/bin/btrfs-manage-raid

